# systemd/weather-daemon.service
[Unit]
Description=Weather Station Data Collector Daemon
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=david
Group=david
WorkingDirectory=/home/david/projects/weatherstation
Environment=PATH=/home/david/projects/weatherstation/venv/bin
ExecStart=/home/david/projects/weatherstation/venv/bin/python main.py --daemon
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Kill the daemon gracefully
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target